# scripts/supervisor/twikifier.conf
[program:twikifier]
command=/usr/local/bin/twikifier
autorestart=true
user=tiddlyweb
redirect_stderr=true
stdout_logfile=syslog
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
