# scripts/twikifier
#/bin/sh
#
# Start-up script for Twikifier
# Credit to https://github.com/chovy/node-startup
#
TS_USER=tiddlyweb
USER_EXEC="su - $TS_USER -c"
NODE_ENV=production
NODE_APP=twikifier
LOGGER="logger -p syslog.info -t twikifier"
SOCKET_FILE=/var/run/twikifier/wst.sock
PID_FILE=/var/run/twikifier.pid

###############

# REDHAT chkconfig header

# chkconfig: - 58 74
# description: twikifier is the script for starting twikifier on boot.
### BEGIN INIT INFO
# Provides: node
# Required-Start:    $network $remote_fs $local_fs 
# Required-Stop:     $network $remote_fs $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: start and stop node
# Description: Node process for app
### END INIT INFO

case "$1" in
    start)
        if [ -f $PID_FILE ]
        then
            echo "$PID_FILE exists, process is already running or crashed"
        else
            echo "Starting twikifier..."
            { NODE_ENV=$NODE_ENV $USER_EXEC $NODE_APP & echo $! > $PID_FILE; } | $LOGGER &
            echo "twikifier started"
        fi
        ;;
	stop)
        if [ ! -f $PID_FILE ]
        then
            echo "$PID_FILE does not exist, process is not running"
        else
            echo "Stopping twikifier ..."
            kill `cat $PID_FILE`
            rm $PID_FILE
            rm $SOCKET_FILE
            echo "twikifier stopped"
        fi
        ;;
    restart)
        if [ ! -f $PID_FILE ]
        then
            echo "$PID_FILE does not exist, process is not running"
        else
            echo "Restarting twikifier ..."
            kill `cat $PID_FILE`
            rm $PID_FILE
            rm $SOCKET_FILE
            sleep 1
            { NODE_ENV=$NODE_ENV $USER_EXEC $NODE_APP & echo $! > $PID_FILE; } | $LOGGER &
            echo "twikifier restarted"
        fi
        ;;
    *)
    echo "Usage: /etc/init.d/twikifier {start|stop|restart}"
    ;;
esac
